# org/ystia/dns/dnsmasq/ansible/playbooks/files/hosts-watcher.path
[Path]
PathChanged=/etc/hosts

[Install]
WantedBy=multi-user.target